# docker-iptables/supervisord.conf
[program:docker-gen]
command=/opt/docker-gen/docker-gen -watch -notify "/usr/bin/perl /tmp/update-iptables.pl" /opt/docker-gen/update-iptables.tmpl /tmp/update-iptables.pl
